# build/supervisord.conf
[supervisord]
nodaemon=true

[program:BIND_named]
command=/usr/sbin/named -f -c /etc/bind/named.conf %(ENV_EXTRA_OPTS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:BIND_dnscrypt-proxy]
command=dnscrypt-proxy -config ./dnscrypt-proxy.toml
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
